# maint/templates/kafka.service.j2
# 示例内容: kafka.service.j2
[Unit]
Description=Kafka Service
